# hn_rerank.timer
[Unit]
Description=Run HN Rerank Update every 4 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
